# scripts/crontab.cron
*/1 * * * * /local/sharelatex/scripts/kill_process_pdf.sh >/dev/null 2>&1
